# etc/systemd/system/mikrotik2pcap.timer
[Unit]
Description=Run microtik2pcap service every minute

[Timer]
OnCalendar=*:0/1
Unit=mikrotik2pcap.service

[Install]
WantedBy=timers.target
